# server-stack/deploy/ubuntu/control-api.service
[Unit]
Description=Deploy control HTTP API (dashboard backend)
After=network-online.target deploy-server.service
Wants=network-online.target

[Service]
Type=simple
User=pirate
Group=pirate
EnvironmentFile=-/etc/pirate-deploy.env
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/control-api --deploy-root /var/lib/pirate/deploy --listen-port 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
