# server-stack/deploy/ubuntu/pirate-host-agent.service
[Unit]
Description=Pirate out-of-band host lifecycle agent (server-stack OTA + reboot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pirate
Group=pirate
EnvironmentFile=-/etc/pirate-host-agent.env
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/pirate-host-agent
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
